# backendMain/documentation/setup_scripts/aeroops-backend.service
[Unit]
Description=AeroOps Backend Django Application
After=network.target postgresql.service

[Service]
Type=exec
User=rajat
Group=rajat
WorkingDirectory=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain
Environment="PATH=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/.venv/bin"
Environment="DJANGO_SETTINGS_MODULE=backendMain.settings"

# Use Gunicorn for production
ExecStart=/bin/bash -c 'cd "/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain" && .venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 --timeout 120 --access-logfile logs/gunicorn-access.log --error-logfile logs/gunicorn-error.log backendMain.wsgi:application'

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aeroops-backend

[Install]
WantedBy=multi-user.target
