# backendMain/docs/setup_scripts/celery-beat.service
[Unit]
Description=Celery Beat Scheduler for AeroOps Backend
After=network.target redis.service

[Service]
Type=simple
User=rajat
Group=rajat
WorkingDirectory=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain
Environment="PATH=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/.venv/bin"
ExecStart=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/.venv/bin/celery -A backendMain beat \
    --loglevel=info \
    --logfile=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/logs/celery-beat.log \
    --pidfile=/tmp/celery-beat.pid \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
