# backendMain/docs/setup_scripts/celery-worker.service
[Unit]
Description=Celery Worker for AeroOps Backend
After=network.target redis.service postgresql.service

[Service]
Type=forking
User=rajat
Group=rajat
WorkingDirectory=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain
Environment="PATH=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/.venv/bin"
ExecStart=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/.venv/bin/celery -A backendMain worker \
    --loglevel=info \
    --logfile=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/logs/celery-worker.log \
    --pidfile=/tmp/celery-worker.pid \
    --detach

ExecStop=/home/rajat/Desktop/AeroOps Intel/aeroScrap_backend/backendMain/.venv/bin/celery -A backendMain control shutdown

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
